# modules/adobe-hbase/templates/service/hbase-zookeeper.erb
#!/bin/bash
# 
# (c) Copyright 2009 Cloudera, Inc.
# 
#   Licensed under the Apache License, Version 2.0 (the "License");
#   you may not use this file except in compliance with the License.
#   You may obtain a copy of the License at
# 
#       http://www.apache.org/licenses/LICENSE-2.0
# 
#   Unless required by applicable law or agreed to in writing, software
#   distributed under the License is distributed on an "AS IS" BASIS,
#   WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
#   See the License for the specific language governing permissions and
#   limitations under the License.
# 
# Starts a Hbase regionserver
# 
# chkconfig: 2345 90 10
# description: Hbase regionserver

#source /etc/rc.d/init.d/functions
#source /home/hadoop/.bashrc 
source /home/<%= user %>/.bashrc 
#RETVAL=0
desc="Hbase zookeeper daemon"

start() {
#  echo -n $"Starting $desc (hbase-zookeeper): "
#  daemon --user <%= user %> <%= hbase_home %>/bin/hbase-daemons.sh --config "<%= hbase_home %>/conf" start zookeeper
	<%= hbase_home %>/bin/hbase-daemons.sh --config "<%= hbase_home %>/conf" start zookeeper
#  RETVAL=$?
#  echo
#  [ $RETVAL -eq 0 ] && touch /var/lock/subsys/hbase-zookeeper
#  return $RETVAL
}

stop() {
#  echo -n $"Stopping $desc (hbase-zookeeper): "
#  daemon --user  <%= user %> <%= hbase_home %>/bin/hbase-daemons.sh --config "<%= hbase_home %>/conf" stop zookeeper
	<%= hbase_home %>/bin/hbase-daemons.sh --config "<%= hbase_home %>/conf" stop zookeeper
#  RETVAL=$?
#  sleep 5
#  echo
#  [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/hbase-zookeeper
}

restart() {
  stop
  start
}

condrestart(){
  [ -e /var/lock/subsys/hbase-zookeeper ] && restart || :
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    restart
    ;;
  condrestart|try-restart)
    condrestart
    ;;
  *)
    echo $"Usage: $0 {start|stop|restart|try-restart}"
    exit 1
esac

exit $RETVAL
